# superivsord.conf
[program:lshs]
command=/home/deploy/websites/jqflask/venv/bin/uwsgi --ini /home/deploy/websites/jqflask/uwsgi.ini